# infra/systemd/example-timers/lumin-agent@09-customer-success.service.d/onboarding-sweep.conf
# Drop-in override for lumin-agent@09-customer-success.service
#
# This demonstrates the systemd drop-in pattern as an ALTERNATIVE approach to
# concrete task service files. A drop-in overrides specific settings from the
# parent template unit (lumin-agent@.service) for this specific instance.
#
# IMPORTANT: Drop-ins affect ALL uses of lumin-agent@09-customer-success.service,
# so this pattern is only suitable when a single timer activates this instance.
# When multiple timers run the same agent with different tasks, use concrete
# named service files (e.g., lumin-agent-09-onboarding-sweep.service) instead.
#
# To use this pattern:
#   sudo mkdir -p /etc/systemd/system/lumin-agent@09-customer-success.service.d/
#   sudo cp onboarding-sweep.conf /etc/systemd/system/lumin-agent@09-customer-success.service.d/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now lumin-agent-09-onboarding-sweep.timer
#
# Then the timer can reference: Unit=lumin-agent@09-customer-success.service

[Service]
Environment=TASK=daily_onboarding_sweep
