# infra/systemd/example-timers/lumin-agent-sbia-followup.timer
[Unit]
Description=Lumin SBIA — Daily Followup Dispatch (14:00 UTC / 10:00 ET)
Documentation=file:/opt/lumin-agents/infra/systemd/README.md

[Timer]
# Daily at 14:00 UTC (10:00 ET) — followup emails/messages to discovered venues
OnCalendar=*-*-* 14:00:00
AccuracySec=1m
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
