# infra/systemd/example-timers/lumin-agent-09-churn-scan.service
[Unit]
Description=Lumin Agent 09 — Daily Churn Scan
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=lumin
Group=lumin
WorkingDirectory=/opt/lumin-agents
EnvironmentFile=/opt/lumin-agents/.env
Environment=TASK=daily_churn_scan

ExecStart=/opt/lumin-agents/agents/agent-09-customer-success/venv/bin/python \
          /opt/lumin-agents/scripts/run_agent.py \
          agent-09-customer-success daily_churn_scan

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/lumin-agents/agents/agent-09-customer-success
TimeoutStartSec=900
Restart=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lumin-agent-09-churn-scan
